# roles/derp/templates/derper.service.j2
[Unit]
Description=Tailscale DERP Server
After=network.target

[Service]
Type=simple
# derper runs HTTPS on internal port 8443
# It auto-generates a self-signed certificate
# Caddy proxies https://relay.xxx -> https://127.0.0.1:8443 with tls_insecure_skip_verify
ExecStart=/opt/derp/derper \
  --hostname={{ tailscale_derp_hostname }} \
  -a=:{{ tailscale_derp_http_port | default(8443) }} \
  --http-port=-1 \
  --stun-port={{ tailscale_derp_stun_port }} \
{% if tailscale_derp_verify_clients %}
  --verify-clients
{% else %}
  --verify-clients=false
{% endif %}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
